# aqi_recoder.timer
[Unit]
Description= aqi recoder timer

[Timer]
OnBootSec=10min
OnUnitActiveSec=15min
Unit=aqi_recoder.service

[Install]
WantedBy=multi-user.target
